# utils/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat8]
directory=/usr/local/tomcat
command=/usr/local/tomcat/bin/catalina.sh run
stdout_logfile=syslog
stderr_logfile=syslog
